# os/packaging/pai-engine/debian/pai-engine.postinst
#!/bin/sh
set -e

case "$1" in
  configure)
    systemctl enable pai-engine.service || true
    ;;
esac

#DEBHELPER#
